# crontab
#run python script every minutes
* * * * * python /hello.py >> /logfile
