# crontab
#run python script every day once
0 0 * * * python /daily_worker.py
